# pkg/dropbear-aa37404/rules/96-dropbear.rules
KERNEL=="loop0", RUN+="/bin/sh -c '/usr/local/dropbear/boot.sh'"
